# utils/rpm/redis-server.target
[Unit]
Description=Redis Server
After=network.target
