# utils/systemd/redis-sentinel.target
[Unit]
Description=Redis Sentinel
# even if redis is an empty target we
# can still launch sentinels.
#
After=redis-server.target
